# cron/webcam.service
[Unit]
Description=Webcam script

[Service]
ExecStart=/home/felix/.dotfiles/cron/webcam.sh
Type=oneshot
