# systemd/borg-notify.timer
[Unit]
Description=10 min borg-notification

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min 

[Install]
WantedBy=timers.target
