# systemd/borg-notify.service
[Unit]
Description=Borg-notification

[Service]
ExecStart=/home/felix/.dotfiles/systemd/borg-notify.sh
Type=oneshot
