# systemd/webcam.service
[Unit]
Description=Webcam script

[Service]
ExecStart=/home/felix/.dotfiles/systemd/webcam.sh
Type=oneshot
